# systemd/vending.service
[Unit]
Description=Water Vending Machine Production Service
After=network.target pigpiod.service
Requires=pigpiod.service

[Service]
Type=simple
User=pi
WorkingDirectory=/opt/water_vending
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/opt/water_vending/deploy/vending_email.env
ExecStart=/opt/water_vending/.venv/bin/python /opt/water_vending/main.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
